# api/supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=app
logfile=/var/log/supervisor/supervisord.log
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:frankenphp]
command=frankenphp run --config /etc/frankenphp/Caddyfile --watch
directory=/app
user=app
autostart=true
autorestart=true
priority=10
stopsignal=SIGTERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:messenger]
command=php /app/bin/console messenger:consume async -vv --time-limit=3600 --memory-limit=256M
directory=/app
user=app
autostart=true
autorestart=true
startretries=999
startsecs=5
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:messenger_outbox]
command=php /app/bin/console messenger:consume outbox -vv --time-limit=3600 --memory-limit=256M
directory=/app
user=app
autostart=true
autorestart=true
startretries=999
startsecs=5
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
